# root/etc/services.d/sstvproxy/run
#!/usr/bin/with-contenv bash

umask 022

cd /app/sstvproxy
exec \
	s6-setuidgid abc /usr/bin/python3 /app/sstvproxy/sstvProxy.py $YAP_ARGS
